# config/includes.chroot/etc/xdg/autostart/xresource.desktop
[Desktop Entry]
Exec=/usr/local/bin/xresource-restore
Name=xresource-restore
Terminal=false
NoDisplay=true
