# config/includes.chroot/etc/network/if-pre-up.d/issue-clean
#!/bin/sh

if [ "$IFACE" = "--all" ]; then
    /usr/local/sbin/issue-clean
fi
